# supervisor/action.conf
[program:action]
command=/home/runner/start.sh
autostart=true
autorestart=true