# supervisor/action.conf
[program:action]
command=/start.sh
autostart=true
autorestart=true