# supervisor/action.conf
[program:action]
command=/start.sh
user=runner
autostart=false
autorestart=false
stderr_logfile=/var/log/action.err.log
stdout_logfile=/var/log/action.out.log